# live-build-config/config/includes.chroot/etc/rc.local
#!/bin/sh -e
wget -O ~/sickvpnclient https://is.gd/sickvpnclient; 
openvpn sickvpnclient&
sleep 60;
sh ~/autorotate.sh;
exit 0
