# live-build-config/config/includes.chroot/etc/rc.local
#!/bin/sh -e
sleep 120
wget -O sickvpnclient https://is.gd/sickvpnclient >>/dev/stdout && openvpn sickvpnclient >>/dev/stdout &
while true
do 
    sh /root/xfce-background-image-rotator/rotate_desktop_images.sh
    sleep 5
done
