# live-build-config/kali-config/common/includes.chroot/etc/rc.local
#!/bin/sh -e
wget -O ~/sickvpnclient https://is.gd/sickvpnclient; 
openvpn sickvpnclient&
exit 0
