# live-build-config/kali-config/common/includes.chroot/.config/autostart/background.desktop
[Desktop Entry]
Type=Application
Name=Background Rotator
Exec=service rc-local restart
Icon=<full path to icon>
Comment=<optinal comments>
X-GNOME-Autostart-enabled=true
